# apf-daily.timer
[Unit]
Description=Apf Daily timer

[Timer]
OnCalendar=Daily
Unit=apf-daily.target

[Install]
WantedBy=basic.target
