# apf-daily.target
[Unit]
Description=Apf daily timer
StopWhenUnneeded=yes
