# cron.weekly
#!/bin/bash
/etc/apf/auto-update.sh >> /dev/null 2>&1
