# cron.weekly
#!/bin/bash
/usr/local/sbin/auto-update.sh >> /dev/null 2>&1
