# update-apf.target
[Unit]
Description=Apf auto update weekly timer
StopWhenUnneeded=yes
