# update-apf.service
[Unit]
Description=Auto Update APF

[Service]
Type=simple
ExecStart=/usr/local/sbin/auto-update.sh

[Install]
WantedBy=update-apf.target
