# pkg/deb/debian/postinst
#!/bin/bash
# apf postinst — generated by pkg_gen.sh from pkg_lib templates
#
# Post-install scriptlet: sets permissions, runs config migration if
# upgrading from install.sh-based installation.
#
#
set -e

LEGACY_PATH="/etc/apf"

case "$1" in
    configure)
        # Set permissions on config files
        command chmod 640 /etc/apf/conf.apf /etc/apf/internals/internals.conf \
            /etc/apf/allow_hosts.rules /etc/apf/deny_hosts.rules \
            /etc/apf/cc_allow.rules /etc/apf/cc_deny.rules \
            /etc/apf/sysctl.rules /etc/apf/sdrop_hosts.rules \
            /etc/apf/silent_ips.rules 2>/dev/null || true  # files may not exist on first install
        command chmod 750 /etc/apf/hook_pre.sh /etc/apf/hook_post.sh 2>/dev/null || true  # hooks may not exist yet
        command chmod 750 /var/lib/apf/tmp /etc/apf/geoip /var/log/apf 2>/dev/null || true  # dirs may not exist yet
        # Auto-detect default network interface on fresh install
        if [ -z "$2" ] && [ ! -f /etc/apf/internals/.apf.restore ]; then
            _detected=$(ip route show default 2>/dev/null | awk '{for(i=1;i<=NF;i++) if($i=="dev"){print $(i+1);exit}}')  # ip may fail in minimal chroot
            if [ -n "$_detected" ] && [ "$_detected" != "eth0" ]; then
                sed -i "s/^IFACE_UNTRUSTED=\"eth0\"/IFACE_UNTRUSTED=\"$_detected\"/" /etc/apf/conf.apf
            fi
        fi
        # Run importconf if migrating from install.sh backup
        if [ -d "${LEGACY_PATH}.bk.last" ]; then
            if [ -x /usr/lib/apf/extras/importconf ]; then
                BK_LAST="${LEGACY_PATH}.bk.last" INSTALL_PATH="$LEGACY_PATH" /usr/lib/apf/extras/importconf || true  # importconf may fail on corrupt backup
            fi
        fi
        # Clean legacy artifacts from install.sh-based installations
        if [ -L "${LEGACY_PATH}.bk.last" ]; then
            # Pre-decomposition files (pre-2.0.2)
            command rm -f "$LEGACY_PATH"/firewall 2>/dev/null || true  # may not exist
            command rm -f "$LEGACY_PATH"/internals/functions.apf 2>/dev/null || true  # may not exist
            command rm -f "$LEGACY_PATH"/internals/geoip.apf 2>/dev/null || true  # may not exist
            command rm -f "$LEGACY_PATH"/internals/ctlimit.apf 2>/dev/null || true  # may not exist
            # Legacy cron entries (replaced by consolidated cron.d.apf)
            command rm -f /etc/cron.hourly/fw 2>/dev/null || true  # may not exist
            command rm -f /etc/cron.daily/fw 2>/dev/null || true  # may not exist
            command rm -f /etc/cron.d/refresh.apf 2>/dev/null || true  # may not exist
            command rm -f /etc/cron.d/apf_develmode 2>/dev/null || true  # may not exist
            command rm -f /etc/cron.d/ctlimit.apf 2>/dev/null || true  # may not exist
        fi
        # Reload systemd
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload 2>/dev/null || true  # systemd may not be running
        fi
        # Disable conflicting firewalls if active
        for _svc in firewalld ufw; do
            if command -v systemctl >/dev/null 2>&1 && systemctl is-active --quiet "$_svc" 2>/dev/null; then  # suppress stderr if service unknown
                systemctl stop "$_svc" 2>/dev/null || true  # service may resist stop
                systemctl mask "$_svc" 2>/dev/null || true  # mask may fail on read-only fs
            fi
        done
        # NOTE: Service is intentionally NOT enabled at boot.
        # Debian Policy 9.3.3.1 recommends against auto-enabling services
        # that modify network access. User must run: systemctl enable apf.service
        #
        # NOTE: VNET generation deferred to first apf -s.
        # vnetgen requires a running network stack and configured interfaces.
        ;;
esac

#DEBHELPER#

exit 0
