# pkg/deb/debian/preinst
#!/bin/bash
# apf preinst — generated by pkg_gen.sh from pkg_lib templates
#
# Pre-install scriptlet: detects and backs up existing install.sh-based
# installations before package lays down FHS files.
#
#
set -e

LEGACY_PATH="/etc/apf"

case "$1" in
    install|upgrade)
        # Detect and back up existing install.sh-based installation
        # A real (non-package) install has the binary as a regular file,
        # not a symlink (which package installs create via symlink farm)
        if [ -f "$LEGACY_PATH/apf" ] && [ ! -L "$LEGACY_PATH/internals/apf.lib.sh" ] && [ ! -L "$LEGACY_PATH/apf" ]; then
            # This is a real (non-package) install — back up
            _bkdir="${LEGACY_PATH}.bk.$(date +%Y%m%d-%s)"
            echo "Backing up existing install.sh installation to $_bkdir"
            command cp -a "$LEGACY_PATH" "$_bkdir"
            command rm -f "${LEGACY_PATH}.bk.last"
            command ln -s "$_bkdir" "${LEGACY_PATH}.bk.last"
            # Preserve state files
            if [ -d "$LEGACY_PATH/tmp" ]; then
                command mkdir -p /var/lib/apf/tmp
                # cp may fail if tmp/ is empty — safe to ignore
                command cp -a "$LEGACY_PATH"/tmp/* /var/lib/apf/tmp/ 2>/dev/null || true
            fi
            # Stop SysVinit services if present — the old install.sh-based
            # installation may have registered init scripts that conflict with
            # the package's own service units; stop and deregister them so the
            # package post-install can set up services cleanly.
            for _initdir in /etc/rc.d/init.d /etc/init.d; do
                if [ -f "$_initdir/apf" ]; then
                    # Service may not be running — safe to ignore stop failure
                    "$_initdir/apf" stop 2>/dev/null || true
                    if command -v chkconfig >/dev/null 2>&1; then
                        # chkconfig may fail if service was never added — safe to ignore
                        chkconfig --del apf 2>/dev/null || true
                    elif command -v update-rc.d >/dev/null 2>&1; then
                        # update-rc.d may fail if links already removed — safe to ignore
                        update-rc.d -f apf remove 2>/dev/null || true
                    fi
                fi
            done
            # Stop systemd services before removing old install — service may
            # not exist or may not be running; either case is safe to ignore
            if command -v systemctl >/dev/null 2>&1; then
                systemctl stop apf.service 2>/dev/null || true
            fi
            # Wipe legacy path so dpkg can extract cleanly. Required because:
            #   1. install.sh creates /etc/apf/{doc,extras} as directories, but
            #      the package ships them as symlinks — dpkg cannot replace dir
            #      with symlink
            #   2. Existing conffiles trigger dpkg conffile prompt; importconf
            #      would then merge old→old, losing new template options
            # Backup at ${LEGACY_PATH}.bk.<ts> preserves all original content;
            # importconf in postinst merges backup values into fresh templates.
            command rm -rf "$LEGACY_PATH" 2>/dev/null || true  # legacy path being wiped after backup
        fi
        ;;
esac

#DEBHELPER#

exit 0
